# roles/avans-entrypoint/templates/etc/systemd/system/certbot-renewal.service.j2
[Unit]
Description=Certbot Renewal

[Service]
ExecStart=/usr/bin/certbot renew --webroot -w "{{ nginx_webroot_directory }}" --post-hook "{{ nginx_host_letsencrypt_post_hook }}"